# alerts/diagnet-alert.cron
*/5 * * * * root /opt/cluster_diag_agent/alerts/diagnet_alert.sh >> /var/log/cluster_diag_agent/alerts_cron.log 2>&1
