# alerts/diagnet-alert.service
[Unit]
Description=DiagNet/VmWatch periodic summary (last 15 minutes)
[Service]
Type=oneshot
ExecStart=/opt/cluster_diag_agent/alerts/diagnet_alert.sh
