# vmwatch/vmwatch.service
[Unit]
Description=VM Watch (synthetic jitter/net/disk)
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
ExecStart=/opt/css_diag_agent/vmwatch/vmwatch.sh
Restart=always
RestartSec=5s
User=root
LimitNOFILE=65536
[Install]
WantedBy=multi-user.target
